# systemd/birdcam.service
[Unit]
Description=BirdCam bird detection system
After=network.target
# Wait for the SSD to be mounted
RequiresMountsFor=/mnt/birdcam

[Service]
Type=simple
User=BIRDCAM_USER
WorkingDirectory=BIRDCAM_DIR
ExecStart=BIRDCAM_DIR/venv/bin/birdcam --config BIRDCAM_DIR/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=2G
CPUQuota=80%

# Ultralytics writes model cache to ~/.config/Ultralytics/
Environment=YOLO_CONFIG_DIR=/mnt/birdcam/.ultralytics

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/birdcam

[Install]
WantedBy=multi-user.target
